# systemd/bible-system-timer.service
[Unit]
Description=AI Bible Documentation System - Periodic Health Check
Documentation=file:///etc/nixos/docs/bibles/
Requires=bible-system.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/etc/nixos
Environment=PYTHONPATH=/etc/nixos/scripts
ExecStart=/etc/nixos/scripts/bible_workflow_manager.py --health-check

# Resource limits for periodic checks
MemoryMax=512M
CPUQuota=25%
TimeoutSec=120

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bible-system-timer

[Install]
WantedBy=timers.target
